# infra/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:terminal-service]
command=/usr/bin/python3 /app/terminal_service.py
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/terminal.err.log
stdout_logfile=/var/log/supervisor/terminal.out.log
environment=OPEN_TERMINAL_API_KEY=%(ENV_OPEN_TERMINAL_API_KEY)s,OPEN_TERMINAL_MULTI_USER=false

[program:little-coder]
command=/usr/local/bin/node /app/little-coder/bin/little-coder.mjs
directory=/app/little-coder
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/little-coder.err.log
stdout_logfile=/var/log/supervisor/little-coder.out.log
environment=OPEN_TERMINAL_URL=http://localhost:8000,OPEN_TERMINAL_API_KEY=%(ENV_OPEN_TERMINAL_API_KEY)s

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
